# pi-setup/instapi-kiosk.service
[Unit]
Description=InstaPi Kiosk Browser
After=instapi.service network.target
Wants=instapi.service

[Service]
Type=simple
User=pi
Environment=DISPLAY=:0
ExecStartPre=/bin/sleep 10
ExecStart=/home/pi/instapi/pi-setup/kiosk.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
